# deploy/systemd/VRX-64-slides.path
[Unit]
Description=Watch the VZGLYD slides directory for package changes

[Path]
PathModified=/var/lib/vzglyd/slides
Unit=VRX-64-slides.service

[Install]
WantedBy=multi-user.target
