# deploy/systemd/vzglyd-slides.path
[Unit]
Description=Watch the VZGLYD slides directory for package changes

[Path]
PathModified=/var/lib/vzglyd/slides
Unit=vzglyd-slides.service

[Install]
WantedBy=multi-user.target
