# deploy/systemd/vzglyd-slides.service
[Unit]
Description=Restart VZGLYD after slides directory changes

[Service]
Type=oneshot
ExecStart=/bin/systemctl try-restart vzglyd.service
